# debian/python3-cyberfusion-proftpd-support.cron.d
*/5 * * * * root /usr/bin/tuxis-cron /etc/cronscripts/proftpd-delete-expired-users.cron
